# rootfs/etc/services.d/crypto-inspect/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Home Assistant Add-on: Crypto Inspect
# Graceful shutdown handler
# ==============================================================================

/bin/true
