# rootfs/etc/s6-overlay/s6-rc.d/crypto-inspect/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Crypto Inspect
# Runs the Crypto Inspect application
# ==============================================================================

bashio::log.info "Starting Crypto Inspect add-on..."

# === Основные настройки с дефолтами ===
export API_PORT=$(bashio::config 'api_port' '9999')
export LOG_LEVEL=$(bashio::config 'log_level' 'info')
export APP_NAME="Crypto Inspect"
export DEBUG=false

# === База данных ===
DB_TYPE=$(bashio::config 'database_type' 'sqlite')
bashio::log.info "Database type: ${DB_TYPE}"

if [[ "${DB_TYPE}" == "postgres" ]]; then
    PG_HOST=$(bashio::config 'postgres_host' '')
    PG_PORT=$(bashio::config 'postgres_port' '5432')
    PG_DB=$(bashio::config 'postgres_database' 'crypto_inspect')
    PG_USER=$(bashio::config 'postgres_username' '')
    PG_PASS=$(bashio::config 'postgres_password' '')

    if [[ -z "${PG_HOST}" ]]; then
        bashio::log.error "PostgreSQL host is required when using postgres database type"
        exit 1
    fi

    export DATABASE_URL="postgresql+asyncpg://${PG_USER}:${PG_PASS}@${PG_HOST}:${PG_PORT}/${PG_DB}"
    export DATABASE_URL_SYNC="postgresql://${PG_USER}:${PG_PASS}@${PG_HOST}:${PG_PORT}/${PG_DB}"
    bashio::log.info "Using PostgreSQL: ${PG_HOST}:${PG_PORT}/${PG_DB}"
else
    # SQLite - встроенная, без зависимостей
    export DATABASE_URL="sqlite+aiosqlite:////data/crypto_inspect.db"
    export DATABASE_URL_SYNC="sqlite:////data/crypto_inspect.db"
    bashio::log.info "Using SQLite: /data/crypto_inspect.db"
fi

# === Redis (опционально) ===
if bashio::config.has_value 'redis_url'; then
    REDIS=$(bashio::config 'redis_url')
    export REDIS_URL="${REDIS}"
    export CELERY_BROKER_URL="${REDIS}"
    export CELERY_RESULT_BACKEND="${REDIS}"
    bashio::log.info "Redis enabled: ${REDIS}"
fi

# === Торговые пары ===
if bashio::config.has_value 'symbols'; then
    # bashio returns array elements on separate lines, join with comma
    export HA_SYMBOLS=$(bashio::config 'symbols' | tr '\n' ',' | sed 's/,$//')
else
    export HA_SYMBOLS="BTCUSDT,ETHUSDT"
fi
bashio::log.info "Symbols: ${HA_SYMBOLS}"

# === Анализ рынка ===
export HA_EXCHANGES="binance,bybit,kraken,kucoin,okx"
export HA_INTERVALS="1m,5m,15m,30m,1h,4h,1d"
export HA_FETCH_INTERVAL=$(bashio::config 'fetch_interval_minutes' '60')
export ANALYSIS_ENABLED=$(bashio::config 'analysis_enabled' 'true')
export ANALYSIS_INTERVAL=$(bashio::config 'analysis_interval_hours' '4')
export ALERT_ON_SIGNALS=$(bashio::config 'alert_on_strong_signals' 'true')
export ALERT_BUY_THRESHOLD=$(bashio::config 'alert_threshold_buy' '75')
export ALERT_SELL_THRESHOLD=$(bashio::config 'alert_threshold_sell' '25')

# === Bybit API ===
if bashio::config.has_value 'bybit_api_key'; then
    export BYBIT_API_KEY=$(bashio::config 'bybit_api_key')
    export BYBIT_API_SECRET=$(bashio::config 'bybit_api_secret' '')
    export BYBIT_TESTNET=$(bashio::config 'bybit_testnet' 'false')
    bashio::log.info "Bybit API enabled"
fi

# === MCP Server ===
export MCP_ENABLED=$(bashio::config 'mcp_enabled' 'false')
export MCP_PORT=$(bashio::config 'mcp_port' '9998')

# === Backfill ===
export BACKFILL_ENABLED=$(bashio::config 'backfill_enabled' 'true')
export BACKFILL_CRYPTO_YEARS=$(bashio::config 'backfill_crypto_years' '10')
export BACKFILL_TRADITIONAL_YEARS=$(bashio::config 'backfill_traditional_years' '1')

# === AI Analysis ===
export AI_ENABLED=$(bashio::config 'ai_enabled' 'false')
if [[ "${AI_ENABLED}" == "true" ]]; then
    export AI_PROVIDER=$(bashio::config 'ai_provider' 'ollama')
    export AI_LANGUAGE=$(bashio::config 'ai_language' 'ru')
    export AI_INTERVAL=$(bashio::config 'ai_analysis_interval_hours' '24')

    if [[ "${AI_PROVIDER}" == "openai" ]]; then
        export OPENAI_API_KEY=$(bashio::config 'openai_api_key' '')
    else
        export OLLAMA_HOST=$(bashio::config 'ollama_host' 'http://localhost:11434')
        export OLLAMA_MODEL=$(bashio::config 'ollama_model' 'llama3.2')
    fi
    bashio::log.info "AI enabled: ${AI_PROVIDER}"
fi

# === Цели и уведомления ===
export GOAL_ENABLED=$(bashio::config 'goal_enabled' 'false')
if [[ "${GOAL_ENABLED}" == "true" ]]; then
    export GOAL_TARGET=$(bashio::config 'goal_target_value' '100000')
    export GOAL_NAME=$(bashio::config 'goal_name' 'Финансовая свобода')
fi

export BRIEFING_MORNING=$(bashio::config 'briefing_morning_enabled' 'true')
export BRIEFING_EVENING=$(bashio::config 'briefing_evening_enabled' 'true')
export NOTIFICATION_MODE=$(bashio::config 'notification_mode' 'smart')

# === Home Assistant secrets ===
export HA_SECRETS_PATH="/config/secrets.yaml"
if [[ -f "${HA_SECRETS_PATH}" ]]; then
    bashio::log.info "secrets.yaml found"
fi

# === Ingress ===
export INGRESS_PATH=$(bashio::addon.ingress_entry)

# === Инициализация ===
mkdir -p /data
mkdir -p /data/huggingface  # Cache for ML models
cd /app

# HuggingFace cache directory for ML models (Chronos, etc.)
export HF_HOME="/data/huggingface"
export TRANSFORMERS_CACHE="/data/huggingface"
export TORCH_HOME="/data/huggingface/torch"

bashio::log.info "Initializing database..."
if [[ "${DB_TYPE}" == "sqlite" ]]; then
    python3 -c "
from models.base import Base
from sqlalchemy import create_engine
engine = create_engine('sqlite:////data/crypto_inspect.db')
Base.metadata.create_all(engine)
print('SQLite tables created')
"
else
    bashio::log.info "Running Alembic migrations..."
    alembic upgrade head || bashio::log.warning "Migration skipped"
fi

bashio::log.info "Starting server on port ${API_PORT}..."
exec python3 -m uvicorn main:app \
    --host 0.0.0.0 \
    --port "${API_PORT}" \
    --log-level "${LOG_LEVEL}"
